# service/explorer.node/run
#!/bin/bash

exec /usr/local/bin/run-service haobtc 'node start.js -r 300 -c bitcoin -s node'
