# service/explorer.txcursor/run
#!/bin/bash

exec /usr/local/bin/run-service haobtc 'mongo --quiet mongoscripts/txcursor.js'
