# service/explorer.query/run
#!/bin/bash

exec /usr/local/bin/run-service haobtc 'node start.js -n 4 -p 9002'
