# service/explorer.mnode/run
#!/bin/bash

exec /usr/local/bin/run-service haobtc 'node start.js -r 300 -c dogecoin -c litecoin -c darkcoin -s mnode'
